# agent/ddns-agent.service
[Unit]
Description=DDNS agent — push current public IP to NodeGet js-worker
Documentation=https://github.com/your-org/nodeget-cloudflare-ddns
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Drop credentials into /etc/ddns-agent/config — root-owned, mode 600.
ExecStart=/usr/local/bin/ddns-agent.sh --verbose
# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/ddns-agent
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6
LockPersonality=true
RestrictNamespaces=true
RestrictSUIDSGID=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
