# agent/ddns-agent.timer
[Unit]
Description=Run ddns-agent every 5 minutes
Documentation=man:systemd.timer(5)

[Timer]
# First run 1 min after boot, then every 5 min.
OnBootSec=1min
OnUnitActiveSec=5min
# Spread load across the fleet if you ever deploy more than one
RandomizedDelaySec=30s
Persistent=true
Unit=ddns-agent.service
AccuracySec=15s

[Install]
WantedBy=timers.target
